# docker/crontab
* * * * * cd /var/www/html && php artisan schedule:run >> /var/log/cron.log 2>&1

# Clean Laravel logs weekly
0 0 * * 0 cd /var/www/html && php artisan log:clean

# Clear expired cache daily
0 0 * * * cd /var/www/html && php artisan cache:clear

# Empty temporary files weekly
0 0 * * 0 cd /var/www/html && php artisan temp:clear

# Leave an empty line at the end of the file
